# packaging/sendry-web.service
[Unit]
Description=Sendry Web - Web Interface for Sendry
Documentation=https://github.com/foxzi/sendry
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sendry
Group=sendry
ExecStart=/usr/bin/sendry-web serve --config /etc/sendry/web.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

# Allow write access to data directories
ReadWritePaths=/var/lib/sendry-web /var/log/sendry

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sendry-web

[Install]
WantedBy=multi-user.target
